# telegram-shell-bot.service
[Unit]
Description=Telegram Shell Bot
After=network.target

[Service]
Type=simple
User=telegram_bot
WorkingDirectory=/home/telegram_bot/telegramShell
# Use unbuffered Python output and redirect all output to journal
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/telegram_bot/telegramShell/venv/bin/python -u bot.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 